# systemd/blackbeard.service
[Unit]
Description=Blackbeard - AI Media Server Manager
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
WorkingDirectory=/home/mermanarchy/blackbeard
ExecStart=/home/mermanarchy/blackbeard/.venv/bin/python -m bot.main
Restart=always
RestartSec=10
User=mermanarchy
EnvironmentFile=/home/mermanarchy/blackbeard/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=blackbeard

[Install]
WantedBy=multi-user.target
